# scripts/com.promptcompactor.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.promptcompactor.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>APFEL_VENV_PYTHON</string>
        <string>COMPACTOR_HOME_PATH/scripts/hook_runner.py</string>
        <string>--serve</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/tmp/apfel-daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/apfel-daemon.stderr.log</string>
</dict>
</plist>
